# freebsd/overlay/usr/local/etc/rc.d/aifw_watchdog
#!/bin/sh
#
# PROVIDE: aifw_watchdog
# REQUIRE: aifw_api
# KEYWORD: shutdown

. /etc/rc.subr

name="aifw_watchdog"
rcvar="aifw_watchdog_enable"

load_rc_config $name

: ${aifw_watchdog_enable:="NO"}
: ${aifw_watchdog_pidfile:="/var/run/aifw_watchdog.pid"}
: ${aifw_watchdog_supervisor_pidfile:="/var/run/aifw_watchdog-supervisor.pid"}
: ${aifw_watchdog_interval:="60"}

pidfile="${aifw_watchdog_supervisor_pidfile}"
procname="/usr/sbin/daemon"
command="/usr/sbin/daemon"
# Watchdog runs as root because it needs to call `service X start`. The
# loop body is in /usr/local/libexec/aifw-watchdog.sh; here we just hand
# it to daemon(8) for supervision.
command_args="-f -p ${aifw_watchdog_pidfile} -P ${aifw_watchdog_supervisor_pidfile} -R 5 -S -T aifw_watchdog -o /var/log/aifw/watchdog.log /usr/local/libexec/aifw-watchdog.sh"

start_precmd="aifw_watchdog_precmd"
stop_cmd="aifw_watchdog_stop"
stop_postcmd="aifw_watchdog_poststop"

aifw_watchdog_precmd()
{
    /bin/mkdir -p /var/log/aifw
    /usr/bin/pkill -f "daemon:.*aifw-watchdog" 2>/dev/null
    /usr/bin/pkill -f "aifw-watchdog.sh" 2>/dev/null
    /bin/rm -f ${aifw_watchdog_pidfile} ${aifw_watchdog_supervisor_pidfile}
    # Export the configured poll interval to the loop script via env.
    export AIFW_WATCHDOG_INTERVAL="${aifw_watchdog_interval}"
}

aifw_watchdog_stop()
{
    # Same SIGTERM/SIGKILL pattern as the other AiFw rc.d scripts so a
    # wedged supervisor can't block the update loop.
    if [ -f "${aifw_watchdog_supervisor_pidfile}" ]; then
        sup_pid=$(cat "${aifw_watchdog_supervisor_pidfile}")
        echo "Stopping aifw_watchdog (supervisor pid ${sup_pid})."
        if kill -TERM "${sup_pid}" 2>/dev/null; then
            i=0
            while [ $i -lt 5 ] && kill -0 "${sup_pid}" 2>/dev/null; do
                sleep 1
                i=$((i+1))
            done
            if kill -0 "${sup_pid}" 2>/dev/null; then
                echo "Graceful stop timed out; sending SIGKILL"
                kill -KILL "${sup_pid}" 2>/dev/null
                /usr/bin/pkill -KILL -f "aifw-watchdog.sh" 2>/dev/null
            fi
        fi
    elif [ -f "${aifw_watchdog_pidfile}" ]; then
        pid=$(cat "${aifw_watchdog_pidfile}")
        kill -TERM "${pid}" 2>/dev/null
        sleep 2
        kill -KILL "${pid}" 2>/dev/null
    else
        echo "aifw_watchdog is not running."
    fi
    /bin/rm -f ${aifw_watchdog_pidfile} ${aifw_watchdog_supervisor_pidfile}
}

aifw_watchdog_poststop()
{
    /bin/rm -f ${aifw_watchdog_pidfile} ${aifw_watchdog_supervisor_pidfile}
}

run_rc_command "$1"
