# freebsd/overlay/usr/local/etc/rc.d/aifw_daemon
#!/bin/sh
#
# PROVIDE: aifw_daemon
# REQUIRE: NETWORKING pf devfs
# KEYWORD: shutdown

. /etc/rc.subr

name="aifw_daemon"
rcvar="aifw_daemon_enable"

load_rc_config $name

: ${aifw_daemon_enable:="NO"}
: ${aifw_daemon_pidfile:="/var/run/aifw_daemon.pid"}
: ${aifw_daemon_supervisor_pidfile:="/var/run/aifw_daemon-supervisor.pid"}
: ${aifw_daemon_db:="/var/db/aifw/aifw.db"}
: ${aifw_daemon_env:=""}

pidfile="${aifw_daemon_supervisor_pidfile}"
procname="/usr/sbin/daemon"
aifw_daemon_binary="/usr/local/sbin/aifw-daemon"
command="/usr/sbin/daemon"
command_args="-f -p ${aifw_daemon_pidfile} -P ${aifw_daemon_supervisor_pidfile} -R 5 -S -T aifw_daemon -o /var/log/aifw/daemon.log -u aifw ${aifw_daemon_binary} --db ${aifw_daemon_db} --log-level info"

start_precmd="aifw_daemon_precmd"
stop_cmd="aifw_daemon_stop"
stop_postcmd="aifw_daemon_poststop"

aifw_demote_if_clustered()
{
    if [ "$(sysrc -n aifw_cluster_enabled 2>/dev/null)" = "YES" ]; then
        sysctl net.inet.carp.demotion=240 >/dev/null 2>&1 || true
        sleep 1
    fi
}

aifw_daemon_precmd()
{
    /bin/mkdir -p /var/log/aifw
    /usr/sbin/chown aifw:aifw /var/log/aifw
    /usr/bin/pkill -f "daemon:.*aifw-daemon" 2>/dev/null
    /usr/bin/pkill -x aifw-daemon 2>/dev/null
    /bin/rm -f ${aifw_daemon_pidfile} ${aifw_daemon_supervisor_pidfile}
    # Pre-create singleton lockfile owned by aifw — /var/run is root-only,
    # so the binary (running as aifw via daemon -u aifw) can't create it.
    /usr/bin/touch /var/run/aifw-daemon.lock
    /usr/sbin/chown aifw:aifw /var/run/aifw-daemon.lock
    # Read the loopback API key from a 640 file (mode root:aifw) rather than
    # from /etc/rc.conf (which is world-readable mode 644).  The file is
    # written by aifw-setup; the daemon process needs it to authenticate
    # background tasks (RoleWatcher, HealthProber, ClusterReplicator) to the
    # local API.
    if [ -r /usr/local/etc/aifw/daemon.key ]; then
        export AIFW_LOOPBACK_API_KEY="$(cat /usr/local/etc/aifw/daemon.key)"
    fi
    # Fallback: honour legacy rc.conf aifw_daemon_env for nodes that have not
    # yet been re-provisioned with the new daemon.key path.
    if [ -z "${AIFW_LOOPBACK_API_KEY}" ] && [ -n "${aifw_daemon_env}" ]; then
        export ${aifw_daemon_env}
    fi
}

aifw_daemon_stop()
{
    aifw_demote_if_clustered
    # SIGTERM with a hard 10-second wall-clock timeout, then SIGKILL.
    # The default rc.subr stop sends SIGTERM and pwait()s indefinitely —
    # if the binary's tokio runtime won't exit (e.g. a background metrics
    # exporter holding the process alive after main returns), service
    # restart wedges forever and blocks the in-product update loop.
    if [ -f "${aifw_daemon_supervisor_pidfile}" ]; then
        sup_pid=$(cat "${aifw_daemon_supervisor_pidfile}")
        echo "Stopping aifw_daemon (supervisor pid ${sup_pid})."
        if kill -TERM "${sup_pid}" 2>/dev/null; then
            i=0
            while [ $i -lt 10 ] && kill -0 "${sup_pid}" 2>/dev/null; do
                sleep 1
                i=$((i+1))
            done
            if kill -0 "${sup_pid}" 2>/dev/null; then
                echo "Graceful stop timed out; sending SIGKILL"
                kill -KILL "${sup_pid}" 2>/dev/null
                /usr/bin/pkill -KILL -x aifw-daemon 2>/dev/null
            fi
        fi
    elif [ -f "${aifw_daemon_pidfile}" ]; then
        pid=$(cat "${aifw_daemon_pidfile}")
        kill -TERM "${pid}" 2>/dev/null
        sleep 2
        kill -KILL "${pid}" 2>/dev/null
    else
        echo "aifw_daemon is not running."
    fi
    /bin/rm -f ${aifw_daemon_pidfile} ${aifw_daemon_supervisor_pidfile}
}

aifw_daemon_poststop()
{
    /bin/rm -f ${aifw_daemon_pidfile} ${aifw_daemon_supervisor_pidfile}
}

run_rc_command "$1"
