# freebsd/overlay/usr/local/etc/rc.d/aifw_api
#!/bin/sh
#
# PROVIDE: aifw_api
# REQUIRE: NETWORKING aifw_daemon aifw_ids
# KEYWORD: shutdown

. /etc/rc.subr

name="aifw_api"
rcvar="aifw_api_enable"

load_rc_config $name

: ${aifw_api_enable:="NO"}
: ${aifw_api_pidfile:="/var/run/aifw_api.pid"}
: ${aifw_api_supervisor_pidfile:="/var/run/aifw_api-supervisor.pid"}
: ${aifw_api_db:="/var/db/aifw/aifw.db"}
: ${aifw_api_listen:="0.0.0.0"}
: ${aifw_api_port:="8080"}

pidfile="${aifw_api_supervisor_pidfile}"
procname="/usr/sbin/daemon"
aifw_api_binary="/usr/local/sbin/aifw-api"
command="/usr/sbin/daemon"
command_args="-f -p ${aifw_api_pidfile} -P ${aifw_api_supervisor_pidfile} -R 5 -S -T aifw_api -o /var/log/aifw/api.log -u aifw ${aifw_api_binary} --db ${aifw_api_db} --listen ${aifw_api_listen}:${aifw_api_port} --ui-dir /usr/local/share/aifw/ui --log-level info"

start_precmd="aifw_api_precmd"
stop_cmd="aifw_api_stop"
stop_postcmd="aifw_api_poststop"

aifw_demote_if_clustered()
{
    if [ "$(sysrc -n aifw_cluster_enabled 2>/dev/null)" = "YES" ]; then
        sysctl net.inet.carp.demotion=240 >/dev/null 2>&1 || true
        sleep 1
    fi
}

aifw_api_precmd()
{
    /bin/mkdir -p /var/log/aifw
    /usr/sbin/chown aifw:aifw /var/log/aifw
    /usr/bin/pkill -f "daemon:.*aifw-api" 2>/dev/null
    /usr/bin/pkill -x aifw-api 2>/dev/null
    /bin/rm -f ${aifw_api_pidfile} ${aifw_api_supervisor_pidfile}
    # Pre-create singleton lockfile owned by aifw — /var/run is root-only,
    # so the binary (running as aifw via daemon -u aifw) can't create it.
    /usr/bin/touch /var/run/aifw-api.lock
    /usr/sbin/chown aifw:aifw /var/run/aifw-api.lock
}

aifw_api_stop()
{
    aifw_demote_if_clustered
    # SIGTERM with a hard 10-second wall-clock timeout, then SIGKILL.
    # The default rc.subr stop sends SIGTERM and pwait()s indefinitely —
    # if the binary's tokio runtime won't exit (e.g. a background metrics
    # exporter holding the process alive after main returns), service
    # restart wedges forever and blocks the in-product update loop.
    if [ -f "${aifw_api_supervisor_pidfile}" ]; then
        sup_pid=$(cat "${aifw_api_supervisor_pidfile}")
        echo "Stopping aifw_api (supervisor pid ${sup_pid})."
        if kill -TERM "${sup_pid}" 2>/dev/null; then
            i=0
            while [ $i -lt 10 ] && kill -0 "${sup_pid}" 2>/dev/null; do
                sleep 1
                i=$((i+1))
            done
            if kill -0 "${sup_pid}" 2>/dev/null; then
                echo "Graceful stop timed out; sending SIGKILL"
                kill -KILL "${sup_pid}" 2>/dev/null
                /usr/bin/pkill -KILL -x aifw-api 2>/dev/null
            fi
        fi
    elif [ -f "${aifw_api_pidfile}" ]; then
        pid=$(cat "${aifw_api_pidfile}")
        kill -TERM "${pid}" 2>/dev/null
        sleep 2
        kill -KILL "${pid}" 2>/dev/null
    else
        echo "aifw_api is not running."
    fi
    /bin/rm -f ${aifw_api_pidfile} ${aifw_api_supervisor_pidfile}
}

aifw_api_poststop()
{
    /bin/rm -f ${aifw_api_pidfile} ${aifw_api_supervisor_pidfile}
}

run_rc_command "$1"
