# freebsd/overlay/usr/local/etc/rc.d/rdhcpd
#!/bin/sh

# PROVIDE: rdhcpd
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="rdhcpd"
rcvar="rdhcpd_enable"

command="/usr/local/sbin/rdhcpd"
command_args="/usr/local/etc/rdhcpd/config.toml"
pidfile="/var/run/${name}.pid"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"
reload_cmd="${name}_reload"
extra_commands="reload"

rdhcpd_start()
{
    if rdhcpd_status >/dev/null 2>&1; then
        echo "${name} is already running."
        return 0
    fi

    mkdir -p /var/db/rdhcpd/leases /var/log/rdhcpd /usr/local/etc/rdhcpd
    chown -R aifw:aifw /var/db/rdhcpd /var/log/rdhcpd /usr/local/etc/rdhcpd

    if [ ! -f /usr/local/etc/rdhcpd/config.toml ]; then
        echo "ERROR: ${name} config not found at /usr/local/etc/rdhcpd/config.toml"
        return 1
    fi

    echo "Starting ${name}."
    /usr/sbin/daemon -p ${pidfile} -o /var/log/rdhcpd/rdhcpd.log ${command} ${command_args}
}

rdhcpd_stop()
{
    if [ -f "${pidfile}" ]; then
        pid=$(cat "${pidfile}")
        echo "Stopping ${name} (pid ${pid})."
        kill "${pid}" 2>/dev/null
        pkill -f "daemon.*rdhcpd" 2>/dev/null
        rm -f "${pidfile}"
        sleep 1
    else
        echo "${name} is not running."
    fi
}

rdhcpd_status()
{
    if [ -f "${pidfile}" ] && kill -0 "$(cat "${pidfile}")" 2>/dev/null; then
        echo "${name} is running (pid $(cat "${pidfile}"))."
        return 0
    else
        echo "${name} is not running."
        return 1
    fi
}

rdhcpd_reload()
{
    if [ -f "${pidfile}" ]; then
        pid=$(cat "${pidfile}")
        echo "Reloading ${name} (pid ${pid})."
        kill -HUP "${pid}" 2>/dev/null
    else
        echo "${name} is not running."
        return 1
    fi
}

load_rc_config $name
: ${rdhcpd_enable:="NO"}
run_rc_command "$1"
