# freebsd/overlay/usr/local/etc/rc.d/rdhcpd
#!/bin/sh
#
# PROVIDE: rdhcpd
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="rdhcpd"
rcvar="rdhcpd_enable"

load_rc_config $name

: ${rdhcpd_enable:="NO"}
: ${rdhcpd_config:="/usr/local/etc/rdhcpd/config.toml"}
: ${rdhcpd_pidfile:="/var/run/rdhcpd/rdhcpd.pid"}
: ${rdhcpd_supervisor_pidfile:="/var/run/rdhcpd/rdhcpd-supervisor.pid"}

pidfile="${rdhcpd_supervisor_pidfile}"
procname="/usr/sbin/daemon"
rdhcpd_binary="/usr/local/sbin/rdhcpd"
command="/usr/sbin/daemon"
command_args="-f -p ${rdhcpd_pidfile} -P ${rdhcpd_supervisor_pidfile} -R 5 -S -T rdhcpd -o /var/log/rdhcpd/rdhcpd.log ${rdhcpd_binary} ${rdhcpd_config}"

start_precmd="rdhcpd_precmd"
stop_postcmd="rdhcpd_poststop"
reload_cmd="rdhcpd_reload"
extra_commands="reload"

rdhcpd_precmd()
{
    /bin/mkdir -p /var/db/rdhcpd/leases /var/log/rdhcpd /usr/local/etc/rdhcpd /var/run/rdhcpd
    /usr/sbin/chown -R aifw:aifw /var/db/rdhcpd /var/log/rdhcpd /usr/local/etc/rdhcpd /var/run/rdhcpd
    /usr/bin/pkill -f "daemon:.*rdhcpd" 2>/dev/null
    /usr/bin/pkill -x rdhcpd 2>/dev/null
    /bin/rm -f ${rdhcpd_pidfile} ${rdhcpd_supervisor_pidfile}

    if [ ! -f ${rdhcpd_config} ]; then
        echo "ERROR: rdhcpd config not found at ${rdhcpd_config}"
        return 1
    fi
}

rdhcpd_poststop()
{
    /bin/rm -f ${rdhcpd_pidfile} ${rdhcpd_supervisor_pidfile}
}

rdhcpd_reload()
{
    if [ -f "${rdhcpd_pidfile}" ]; then
        pid=$(cat "${rdhcpd_pidfile}")
        echo "Reloading rdhcpd (pid ${pid})."
        /bin/kill -HUP "${pid}" 2>/dev/null
    else
        echo "rdhcpd is not running."
        return 1
    fi
}

run_rc_command "$1"
