# freebsd/overlay/usr/local/etc/rc.d/rdns
#!/bin/sh
#
# PROVIDE: rdns
# REQUIRE: NETWORKING
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable rdns:
#   rdns_enable="YES"

. /etc/rc.subr

name="rdns"
rcvar="rdns_enable"

load_rc_config $name

: ${rdns_enable:="NO"}
: ${rdns_config:="/usr/local/etc/rdns/rdns.toml"}
: ${rdns_pidfile:="/var/run/rdns/rdns.pid"}
: ${rdns_supervisor_pidfile:="/var/run/rdns/rdns-supervisor.pid"}

# Track the daemon(8) supervisor so `service rdns stop` kills both it and
# the rdns child. Using -p alone (child pidfile) leaves the supervisor
# running, and -R 5 respawns the child — accumulating orphans on every
# restart.
pidfile="${rdns_supervisor_pidfile}"
# procname matches the process whose PID is in pidfile — that's daemon(8),
# not rdns itself. rc.subr uses this to detect stale pidfiles.
procname="/usr/sbin/daemon"
rdns_binary="/usr/local/sbin/rdns"
command="/usr/sbin/daemon"
command_args="-f -p ${rdns_pidfile} -P ${rdns_supervisor_pidfile} -R 5 -S -T rdns -o /var/log/rdns/rdns.log ${rdns_binary} -c ${rdns_config}"

start_cmd="rdns_start"
stop_postcmd="rdns_poststop"

rdns_start()
{
    /bin/mkdir -p /var/run/rdns /usr/local/etc/rdns/zones /usr/local/etc/rdns/rpz
    # aifw user needs write to /usr/local/etc/rdns/rpz to stage blocklist
    # files. Keeping the group as `rdns` so rDNS itself keeps reading them.
    /usr/sbin/chown -R aifw:rdns /usr/local/etc/rdns/rpz 2>/dev/null
    /bin/chmod 2775 /usr/local/etc/rdns/rpz 2>/dev/null
    # Reap any orphan daemon(8)+rdns pairs left by previous buggy restarts.
    /usr/bin/pkill -f "daemon:.*${procname}" 2>/dev/null
    /usr/bin/pkill -x rdns 2>/dev/null
    /bin/rm -f ${rdns_pidfile} ${rdns_supervisor_pidfile}

    echo "Starting rdns."
    ${command} ${command_args}

    # Wait up to 3 s for the control socket to appear, then make it writable
    # by the aifw group so aifw-api / aifw-daemon can send stats-json /
    # tail-blocks / reload-rpz commands. On FreeBSD, connect(2) on a Unix
    # domain socket requires write perm on the socket inode.
    i=0
    while [ ! -S /var/run/rdns/control.sock ] && [ $i -lt 30 ]; do
        /bin/sleep 0.1
        i=$((i+1))
    done
    if [ -S /var/run/rdns/control.sock ]; then
        /usr/sbin/chown root:aifw /var/run/rdns/control.sock 2>/dev/null
        /bin/chmod 660 /var/run/rdns/control.sock 2>/dev/null
    fi
}

rdns_poststop()
{
    /bin/rm -f ${rdns_pidfile} ${rdns_supervisor_pidfile}
}

run_rc_command "$1"
