# freebsd/overlay/usr/local/etc/rc.d/trafficcop
#!/bin/sh
#
# PROVIDE: trafficcop
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="trafficcop"
rcvar="trafficcop_enable"

load_rc_config $name

: ${trafficcop_enable:="NO"}
: ${trafficcop_config:="/usr/local/etc/trafficcop/config.yaml"}
: ${trafficcop_pidfile:="/var/run/trafficcop/trafficcop.pid"}
: ${trafficcop_supervisor_pidfile:="/var/run/trafficcop/trafficcop-supervisor.pid"}

pidfile="${trafficcop_supervisor_pidfile}"
procname="/usr/sbin/daemon"
trafficcop_binary="/usr/local/sbin/trafficcop"
command="/usr/sbin/daemon"
command_args="-f -p ${trafficcop_pidfile} -P ${trafficcop_supervisor_pidfile} -R 5 -S -T trafficcop -o /var/log/trafficcop/trafficcop.log -u aifw ${trafficcop_binary} -c ${trafficcop_config}"

start_precmd="trafficcop_precmd"
stop_postcmd="trafficcop_poststop"

trafficcop_precmd()
{
    /bin/mkdir -p /var/run/trafficcop /var/log/trafficcop
    /usr/sbin/chown aifw:aifw /var/log/trafficcop /var/run/trafficcop
    /usr/bin/pkill -f "daemon:.*trafficcop" 2>/dev/null
    /usr/bin/pkill -x trafficcop 2>/dev/null
    /bin/rm -f ${trafficcop_pidfile} ${trafficcop_supervisor_pidfile}
    # Pre-create singleton lockfile owned by aifw — trafficcop runs as aifw
    # (via daemon -u aifw) and cannot create files in /var/run/ itself.
    /usr/bin/touch /var/run/trafficcop.lock
    /usr/sbin/chown aifw:aifw /var/run/trafficcop.lock

    if [ -f ${trafficcop_config} ]; then
        if ! ${trafficcop_binary} -c ${trafficcop_config} --validate >/dev/null 2>&1; then
            echo "ERROR: trafficcop config validation failed"
            return 1
        fi
    fi
}

trafficcop_poststop()
{
    /bin/rm -f ${trafficcop_pidfile} ${trafficcop_supervisor_pidfile}
}

run_rc_command "$1"
