# freebsd/overlay/usr/local/etc/rc.d/trafficcop
#!/bin/sh
#
# PROVIDE: trafficcop
# REQUIRE: NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="trafficcop"
rcvar="trafficcop_enable"

load_rc_config $name

: ${trafficcop_enable:="NO"}
: ${trafficcop_config:="/usr/local/etc/trafficcop/config.yaml"}
: ${trafficcop_pidfile:="/var/run/trafficcop/trafficcop.pid"}
: ${trafficcop_supervisor_pidfile:="/var/run/trafficcop/trafficcop-supervisor.pid"}

pidfile="${trafficcop_supervisor_pidfile}"
procname="/usr/sbin/daemon"
trafficcop_binary="/usr/local/sbin/trafficcop"
command="/usr/sbin/daemon"
command_args="-f -p ${trafficcop_pidfile} -P ${trafficcop_supervisor_pidfile} -R 5 -S -T trafficcop -o /var/log/trafficcop/trafficcop.log -u aifw ${trafficcop_binary} -c ${trafficcop_config}"

start_precmd="trafficcop_precmd"
stop_cmd="trafficcop_stop"
stop_postcmd="trafficcop_poststop"

trafficcop_precmd()
{
    /bin/mkdir -p /var/run/trafficcop /var/log/trafficcop
    /usr/sbin/chown aifw:aifw /var/log/trafficcop /var/run/trafficcop
    /usr/bin/pkill -f "daemon:.*trafficcop" 2>/dev/null
    /usr/bin/pkill -x trafficcop 2>/dev/null
    /bin/rm -f ${trafficcop_pidfile} ${trafficcop_supervisor_pidfile}
    # Pre-create singleton lockfile owned by aifw — trafficcop runs as aifw
    # (via daemon -u aifw) and cannot create files in /var/run/ itself.
    /usr/bin/touch /var/run/trafficcop.lock
    /usr/sbin/chown aifw:aifw /var/run/trafficcop.lock

    if [ -f ${trafficcop_config} ]; then
        if ! ${trafficcop_binary} -c ${trafficcop_config} --validate >/dev/null 2>&1; then
            echo "ERROR: trafficcop config validation failed"
            return 1
        fi
    fi
}

trafficcop_stop()
{
    # SIGTERM with a hard 10-second wall-clock timeout, then SIGKILL.
    # The default rc.subr stop sends SIGTERM and pwait()s indefinitely —
    # if the binary's tokio runtime won't exit (e.g. a background metrics
    # exporter holding the process alive after main returns), service
    # restart wedges forever and blocks the in-product update loop.
    if [ -f "${trafficcop_supervisor_pidfile}" ]; then
        sup_pid=$(cat "${trafficcop_supervisor_pidfile}")
        echo "Stopping trafficcop (supervisor pid ${sup_pid})."
        if kill -TERM "${sup_pid}" 2>/dev/null; then
            i=0
            while [ $i -lt 10 ] && kill -0 "${sup_pid}" 2>/dev/null; do
                sleep 1
                i=$((i+1))
            done
            if kill -0 "${sup_pid}" 2>/dev/null; then
                echo "Graceful stop timed out; sending SIGKILL"
                kill -KILL "${sup_pid}" 2>/dev/null
                /usr/bin/pkill -KILL -x trafficcop 2>/dev/null
            fi
        fi
    elif [ -f "${trafficcop_pidfile}" ]; then
        pid=$(cat "${trafficcop_pidfile}")
        kill -TERM "${pid}" 2>/dev/null
        sleep 2
        kill -KILL "${pid}" 2>/dev/null
    else
        echo "trafficcop is not running."
    fi
    /bin/rm -f ${trafficcop_pidfile} ${trafficcop_supervisor_pidfile}
}

trafficcop_poststop()
{
    /bin/rm -f ${trafficcop_pidfile} ${trafficcop_supervisor_pidfile}
}

run_rc_command "$1"
